# contrib/bitbackup.timer
[Unit]
Description=bitbackup-timer

[Timer]
OnCalendar=daily

[X-Fleet]
X-ConditionMachineOf=bitbackup.service